# 构建deb包实例代码/hello_deb/etc/systemd/system/hello.service
[Unit]

Description = hello daemon

[Service]
ExecStart = /opt/hello_deb/hello_deb.sh
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target


