# 构建deb包实例代码/hello_deb/DEBIAN/postinst
#!/bin/bash

echo "......................................................"
echo "..............hello_deb installing...................."
echo "......................................................"

systemctl enable hello
systemctl start hello 


